# demo/deno_php82/supervisor-services.conf
[program:php-fpm]
command=php-fpm -F -R
autostart=true
autorestart=true
startretries=3
startsecs=2
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Run as root so it can manage children (official image default)
user=root

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
startsecs=1
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root

[program:deno]
command=/usr/local/bin/deno-app
autostart=true
autorestart=true
startretries=5
startsecs=2
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
environment=PORT="8000"
